# contrib/systemd-system/vault-certbot.service
[Unit]
Description=Store SSL certificates generated by certbot in the HashiCorp Vault
After=network-online.target
Wants=network-online.target
After=vault.service
Wants=vault.service

[Service]
Type=exec
User=root
Group=root
Environment=PYTHONDONTWRITEBYTECODE=1
WorkingDirectory=/opt/vault-certbot
ExecStart=/opt/vault-certbot/vault-certbot service
TimeoutStopSec=20s
Restart=always
RestartSec=60s

[Install]
WantedBy=multi-user.target
